# rc.local
sh /home/adminw/.___/script/startmotion
exit 0  
